# root/etc/crontabs/crontab
# min   hour    day     month   weekday command
*/5     *       *       *       *       /update-route53-record.sh 2>&1
